# deploy/excalidraw-agg.service
[Unit]
Description=Excalidraw Aggregation Server
After=network.target

[Service]
Type=simple
User=naze
WorkingDirectory=/opt/apps/excalidraw-agg/deploy
Environment=NODE_ENV=production
Environment=FRONTEND_PORT=18790
Environment=STORAGE_DIR=/opt/apps/excalidraw-agg/data
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=excalidraw-agg

[Install]
WantedBy=multi-user.target
